# voice-retro.service
[Unit]
Description=Voice Retro Bot - Telegram Voice Retrospective Bot (Polling Mode)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=voice-retro
Group=voice-retro
WorkingDirectory=/opt/voice-retro
Environment=PATH=/opt/voice-retro/venv/bin
ExecStart=/opt/voice-retro/venv/bin/python local_bot.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
KillMode=mixed
TimeoutStopSec=10

# Polling-specific settings
StartLimitInterval=60
StartLimitBurst=3

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/voice-retro
ProtectHome=true

# Resource limits
LimitNOFILE=65536

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=voice-retro

[Install]
WantedBy=multi-user.target